# pi-db/foyer-db.service
[Unit]
Description=Foyer DB host (SQLite over HTTP)
After=network.target
# Database files live on the USB drive — don't start until it's mounted.
RequiresMountsFor=/mnt/foyerdb

[Service]
# Must match the DB_HTTP_SECRET Cloudflare Pages secret.
Environment=DB_SECRET=dgYMH4ZI5DzazSeX1yw18ULJYHTniyHcSRX-Bmum-Ag
# Data on the USB (ext4). Code stays on the SD so the service can always start.
Environment=DATA_DIR=/mnt/foyerdb/data
Environment=PORT=8787
Environment=ALLOWED_DBS=lanson,burzer
# If the USB is FAT/exFAT instead of ext4, uncomment the next line (WAL needs ext4):
#Environment=JOURNAL_MODE=DELETE
ExecStart=/usr/bin/python3 /home/zoop/foyer-db/server.py
Restart=always
RestartSec=2
User=zoop

[Install]
WantedBy=multi-user.target
